# deploy/systemd/zbarber-api-test.service
[Unit]
Description=Zbarber API (test)
After=network.target

[Service]
Type=simple
User=zbarber
Group=zbarber
WorkingDirectory=/srv/zbarber/test/artifacts/api-server
EnvironmentFile=/etc/zbarber/api-test.env
ExecStart=/usr/bin/node --enable-source-maps dist/index.mjs
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
